# Crypt/src/crypt.service
[[Unit]
Description=B
After=network.target
[Service]
User=client
WorkingDirectory=/home/client/Abhay/cryptage/Crypt/src/
ExecStart=/home/client/Abhay/cryptage/Crypt/src/crypt.sh start

[Install]
WantedBy=multi-user.target
