# Crypt/src/crypt.service
[Unit]
Description=Data encryption
After=network-online.target
[Service]
ExecStart=/bin/crypt.sh
[Install]
WantedBy=multi-user.target
