# Crypt/src/crypt.service
[Unit]
Description=Data encryption
After=network-online.target
[Service]
WorkingDirectory=/home/client/Abhay/cryptage/Crypt/src
ExecStart=/home/client/Abhay/cryptage/Crypt/src/crypt.sh
[Install]
WantedBy=multi-user.target
