# Crypt/src/crypt.service
[Unit]
Description=Data encryption
After=network-online.target

[Service]
User=client
WorkingDirectory=~/Abhay/cryptage/
ExecStart=~/Abhay/cryptage/crypt.sh start

[Install]
WantedBy=multi-user.target
