# Crypt/src/crypt.service
[Unit]
Description=Data encryption
After=network-online.target
[Service]
Type-simple
ExecStart=/home/client/Desktop/Crypt/src/crypt.sh
[Install]
WantedBy=multi-user.target
