# Crypt/src/crypt.service
[Unit]
Description=Data encryption
After=network-online.target

[Service]
User=client
WorkingDirectory=/home/$User/Abhay/cryptage/Crypt/src
ExecStart=/home/client/Abhay/cryptage/Crypt/src/crypt.sh start

[Install]
WantedBy=multi-user.target
