# Crypt/src/decrypt.service
[Unit]
Description=Data decryption
After=network-online.target

[Service]
User=client
WorkingDirectory=/home/$User/Abhay/cryptage/Crypt/src
ExecStart=/home/$User/Abhay/cryptage/Crypt/src/decrypt.sh start

[Install]
WantedBy=multi-user.target
