# Crypt/src/decrypt.service
[Unit]
Description=Decrypt data
After=network.target

[Service]
WorkingDirectory=/home/client/Abhay/cryptage/Crypt/src/
ExecStart=/home/client/Abhay/cryptage/Crypt/src/decrypt.sh start

[Install]
WantedBy=multi-user.target
