# Crypt/src/decrypt.service
[Unit]
Description=Data decryption
After=network-online.target
[Service]
ExecStart=/bin/decrypt.sh
[Install]
WantedBy=multi-user.target
